# assets/ssl_renew.service
[Unit]
Description="Service to run orangehrm ssl renew"

[Service]
ExecStart=/usr/bin/orangehrm ssl renew